# base/src/config/program-rsyslog.conf
[program:rsyslog]
command = /etc/init.d/rsyslog start
