# base/src/config/program-crond.conf
[program:crond]
command = /etc/init.d/crond start
