# pgsql/program-postgresql.conf
[program:postgresql]
command = /bin/bash /opt/startup.sh
